# centos/etc/rc.local
XSym
0013
54c566ec88122e82b24874941e63d0c5
rc.d/rc.local
                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                  